# system/mmdvm.service
[Unit]
Description=MMDVMHost
After=systemd-user-session.service qnrelay.service

[Service]
Type=simple
ExecStart=/usr/local/bin/MMDVMHost /usr/local/etc/MMDVM.qn
Restart=always

[Install]
WantedBy=multi-user.target
