# system/mmdvm.timer
[Timer]
OnStartupSec=30

[Install]
WantedBy=multi-user.target

